# amule/root/etc/s6-overlay/s6-rc.d/svc-amule/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Caught SIGTERM signal!"

HOME=${CONFIGDIR} amulecmd --command shutdown
